# Task_1_Deploy_Systemd_Service/simple-http-server.service
[Unit]
Description=Simple HTTP Server Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/simple-http-server
ExecStart=/usr/bin/python3 /opt/simple-http-server/simple-http-server.py
StandardOutput=append:/var/log/simple-http-server/output.log
StandardError=append:/var/log/simple-http-server/error.log
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target